# deploy/lightsail/templates/systemd/aegis-snapshot.service.tmpl
[Unit]
Description=Create daily Lightsail snapshot for Aegis instance
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aegis_snapshot.sh
